# deploy/supervisor/laravel-scheduler.conf
; Laravel Scheduler — keeps `php artisan schedule:work` running as a daemon.
; schedule:work fires schedule:run every minute internally (no cron needed).
;
; INSTALL ON SERVER:
;   sudo cp /var/www/html/deploy/supervisor/laravel-scheduler.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start laravel-scheduler
;
; Replace /var/www/html with your actual deployment path.
; Replace www-data with the user that owns your Laravel files (check: ls -la /var/www/html).

[program:laravel-scheduler]
process_name=%(program_name)s
command=php /var/www/html/artisan schedule:work
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=30
